# backend/deploy/systemd/node_exporter.service
[Unit]
Description=Prometheus node_exporter (system metrics)
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=node_exporter
Group=node_exporter
# Bind to localhost only — scraped by the local Prometheus, never exposed publicly.
ExecStart=/opt/monitoring/node_exporter/node_exporter \
    --web.listen-address=127.0.0.1:9100
Restart=always
RestartSec=3
# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
